# deploy/fera-gateway.service
[Unit]
Description=Fera Gateway
After=network.target fera-memory.service
Requires=fera-memory.service

[Service]
Type=simple
User=fera
Environment=FERA_HOME=/home/fera
LoadCredential=CLAUDE_CODE_OAUTH_TOKEN:/etc/fera/credentials/claude-code-oauth-token
LoadCredential=TELEGRAM_BOT_TOKEN:/etc/fera/credentials/telegram-bot-token
LoadCredential=TELEGRAM_ALLOWED_USER_ID:/etc/fera/credentials/telegram-allowed-user-id
LoadCredential=MATTERMOST_MAIN_TOKEN:/etc/fera/credentials/mattermost-main-token
ExecStart=/opt/fera-venv/bin/fera-gateway
Restart=on-failure
RestartSec=5
TimeoutStopSec=45
MemoryMax=4G
MemorySwapMax=0

[Install]
WantedBy=multi-user.target
