# deploy/fera-webui.service
[Unit]
Description=Fera Web UI
After=network.target fera-gateway.service

[Service]
Type=simple
User=fera
Environment=FERA_HOME=/home/fera
EnvironmentFile=-/etc/fera/fera-webui.env
ExecStart=/opt/fera-venv/bin/fera-webui
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
